# package/DEBIAN/postinst
sed -i -e 's/enabled=1/enabled=0/' /etc/default/apport
systemctl restart apport

# Deshabilita las actualizaciones automáticas
cat <<EOF >/etc/apt/apt.conf.d/10periodic
APT::Periodic::Update-Package-Lists "0";
APT::Periodic::Download-Upgradeable-Packages "0";
APT::Periodic::AutocleanInterval "0";
APT::Periodic::Unattended-Upgrade "0";
EOF

# Configura el servidor de VNC

cat <<EOF >/etc/systemd/system/x11vnc.service
[Unit]
Description="x11vnc"
Requires=display-manager.service
After=display-manager.service

[Service]
ExecStart=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :0 -auth guess -rfbauth /etc/x11vnc.pass
ExecStop=/usr/bin/killall x11vnc
Restart=on-failure
Restart-sec=2

[Install]
WantedBy=multi-user.target
EOF

systemctl daemon-reload
systemctl start x11vnc
systemctl enable x11vnc

# Configura GRUB sin tiempo de espera

sed -i -e 's/GRUB_TIMEOUT=0/GRUB_TIMEOUT=-1/' /etc/default/grub
update-grub

