# package/DEBIAN/postinst
# Establece el inicio de sesión automático con el usuario por defecto; 'usuario'

echo "[SeatDefaults]\ngreeter-session=lightdm-gtk-greeter\nautologin-user=usuario\n" > /usr/share/lightdm/lightdm.conf.d/60-lightdm-gtk-greeter.conf

# Activa el inicio automático de Zoom al iniciar sesión

cat <<EOF >/home/usuario/.config/autostart/zoom.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/zoom
Hidden=false
X-MATE-Autostart-enabled=true
Name[gl_ES]=Zoom
Name=Zoom
Comment[gl_ES]=Launch Zoom
Comment=Launch Zoom
EOF

# Configura el servidor de VNC
# Hay que crear de forma manual el fichero con la contraseña de acceso 'x11vnc.pass'

cat <<EOF >/etc/systemd/system/x11vnc.service
[Unit]
Description="x11vnc"
Requires=display-manager.service
After=display-manager.service

[Service]
ExecStart=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :0 -auth guess -rfbauth /etc/x11vnc.pass
ExecStop=/usr/bin/killall x11vnc
Restart=on-failure
Restart-sec=2

[Install]
WantedBy=multi-user.target
EOF

systemctl daemon-reload
systemctl enable --now x11vnc
