# package/etc/xdg/autostart/guess.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/guess
Hidden=False
X-MATE-autostart-enabled=true
Name[ES_es]=Guess message
Name=Guess message
